# lcd-cast.service
[Unit]
Description=LCD Desktop Cast Service
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/lcd/Python
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/lcd/Python/main.py
Restart=always
RestartSec=2
KillMode=process
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lcd-cast

[Install]
WantedBy=multi-user.target
