# lcd-cast-user.service
[Unit]
Description=LCD Desktop Cast Service (per-user)
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/home/pi/lcd/Python
# Load GUI-related environment variables from the user's env file (created by installer)
EnvironmentFile=%h/.config/lcd-cast.env
ExecStart=/usr/bin/python3 /home/pi/lcd/Python/main.py
Restart=always
RestartSec=5
# For per-user services the unit must not set User=; it runs as the owning user
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
