# lcd-cast.service
[Unit]
Description=LCD Desktop Cast Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/lcd/Python
Environment=PYTHONUNBUFFERED=1
# Run as root so the service can access GPIO, SPI, I2C, and the LCD device
# nodes without relying on per-user device permissions.
# If you are using the desktop-capture variant, prefer the per-user unit.
EnvironmentFile=-/home/pi/.config/lcd-cast.env
# Fallback defaults if the env file is missing
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/python3 /home/pi/lcd/Python/main.py
Restart=always
RestartSec=2
KillMode=process
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lcd-cast

[Install]
WantedBy=multi-user.target
